# services/dht-search-api/systemd/dht-api.service
[Unit]
Description=DHT Search API Server
Documentation=https://github.com/profullstack/music-torrent
After=network.target bitmagnet.service

[Service]
Type=simple
User=dht
Group=dht
WorkingDirectory=/opt/dht-api

# Environment variables - loaded from file
EnvironmentFile=/opt/dht-api/.env

# Start the API server with pnpm/tsx
ExecStart=/bin/bash -c 'source ~/.bashrc && pnpm start'
Environment=PATH=/home/dht/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/dht-api.log
StandardError=append:/var/log/dht-api.error.log

# Resource limits
LimitNOFILE=65535
MemoryMax=512M

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/dht-api /var/log

[Install]
WantedBy=multi-user.target
